# Paxos_2_1/tgt/DEBIAN.tgtd.css/prerm
#!/bin/sh

# rm -f /etc/init.d/cssd
# update-rc.d cssd remove
