# PaxosProducts/PaxosMemcache/service/upstart/paxos-memcache.conf
# PAXOS memcache

description     "PAXOS memcache"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid paxos

expect fork
exec /home/paxos/bin/paxos memcache start
